# src/rootfs/etc/services.d/vsftpd/run
#!/bin/sh
set -e

s6-setuidgid root vsftpd /etc/vsftpd/vsftpd.conf
